# systemd/rag-backend.service
[Unit]
Description=RAG Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rag-app/backend
Environment="PATH=/usr/bin:/home/ubuntu/rag-app/backend/venv/bin"
ExecStart=/home/ubuntu/rag-app/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rag-backend

[Install]
WantedBy=multi-user.target
